# config/cron.d/distant-reader
#
# Periodic job to scrape the public and private carrels into the patron
# database for the web front-end
#
# Since the app user does not have a home directory, running this as
# me (dbrower). This is not ideal.
#
MAILTO=dbrower
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
SHELL=/bin/bash

*/5 * * * * dbrower /opt/reader/scan_carrels.sh /data-disk/www/html/library/patrons/
# only run this every 15 minutes, and offset from the other script by 3 minutes
3,18,33,48 * * * * dbrower /opt/reader/scan_carrels.sh /data-disk/www/html/library/carrels/ public
